# .config/systemd/user/compton.service
Description=compton
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/compton

[Install]
WantedBy=graphical-session.target
