# .config/systemd/user/user-graphical-session.target
[Unit]
Description=User Graphical Session
Requires=default.target
After=default.target
